# rules/DE.AE/DE.AE-02_windows-endpoint-batch2.xml
<!-- 
  Batch 2: Windows Endpoint Depth (Rules 100027 - 100034)
  NIST Focus: DE.CM-09, DE.AE-02, PR.AA-01, RS.MA-02
-->

<group name="windows, sysmon, security, nist_csf, mitre_attack,">

  <!-- Rule 100027: PowerShell Encoded Command Execution -->
  <rule id="100027" level="10">
    <decoded_as>json</decoded_as>
    <field name="win.system.providerName">Microsoft-Windows-Sysmon</field>
    <field name="win.system.eventID">^1$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(powershell|pwsh)\.exe.*(\-e(nc(oded(Command)?)?)?|\-e\s)</field>
    <description>PowerShell encoded command execution (Defense Evasion)</description>
    <mitre>
      <id>T1059.001</id>
    </mitre>
    <group>DE.AE-02</group>
  </rule>

  <!-- Rule 100028: PowerShell Download Cradle -->
  <rule id="100028" level="10">
    <decoded_as>json</decoded_as>
    <field name="win.system.providerName">Microsoft-Windows-Sysmon</field>
    <field name="win.system.eventID">^1$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(Net\.WebClient|Invoke-WebRequest|DownloadString|DownloadFile|IEX|Invoke-Expression)</field>
    <description>PowerShell web download string executed (Execution/C2)</description>
    <mitre>
      <id>T1059.001</id>
    </mitre>
    <group>DE.AE-02</group>
  </rule>

  <!-- Rule 100029: Scheduled Task Created via Schtasks.exe -->
  <rule id="100029" level="8">
    <decoded_as>json</decoded_as>
    <field name="win.system.providerName">Microsoft-Windows-Sysmon</field>
    <field name="win.system.eventID">^1$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)schtasks\.exe.*\/create</field>
    <description>Scheduled task creation detected via Sysmon (Persistence)</description>
    <mitre>
      <id>T1053.005</id>
    </mitre>
    <group>DE.AE-02</group>
  </rule>

  <!-- Rule 100030: Registry Run/RunOnce Key Modification -->
  <rule id="100030" level="10">
    <decoded_as>json</decoded_as>
    <field name="win.system.providerName">Microsoft-Windows-Sysmon</field>
    <field name="win.system.eventID">^12$|^13$|^14$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)(Software\\Microsoft\\Windows\\CurrentVersion\\Run|Software\\Microsoft\\Windows\\CurrentVersion\\RunOnce)</field>
    <description>Registry Run/RunOnce key modification (Persistence)</description>
    <mitre>
      <id>T1547.001</id>
    </mitre>
    <group>DE.AE-02</group>
  </rule>

  <!-- Rule 100031: Volume Shadow Copy Deleted -->
  <rule id="100031" level="14">
    <decoded_as>json</decoded_as>
    <field name="win.system.providerName">Microsoft-Windows-Sysmon</field>
    <field name="win.system.eventID">^1$</field>
    <field name="win.eventdata.commandLine" type="pcre2">(?i)(vssadmin|wmic).*(delete\s+shadows|shadowcopy\s+delete)</field>
    <description>Volume Shadow Copy deletion detected (Impact)</description>
    <mitre>
      <id>T1490</id>
    </mitre>
    <group>DE.AE-02</group>
  </rule>

  <!-- Rule 100032: New Local Administrator Account Created -->
  <rule id="100032" level="10">
    <decoded_as>json</decoded_as>
    <field name="win.system.providerName">Microsoft-Windows-Security-Auditing</field>
    <field name="win.system.eventID">^4732$|^4728$</field>
    <field name="win.eventdata.targetUserName" type="pcre2">(?i)(Administrators|Domain\s+Admins|Enterprise\s+Admins)</field>
    <description>User added to privileged administrative group (Privilege Escalation)</description>
    <mitre>
      <id>T1078</id>
    </mitre>
    <group>PR.AA-01, PR.AC-04</group>
  </rule>

  <!-- Rule 100033: RDP Enabled or Port Changed -->
  <rule id="100033" level="9">
    <decoded_as>json</decoded_as>
    <field name="win.system.providerName">Microsoft-Windows-Sysmon</field>
    <field name="win.system.eventID">^12$|^13$|^14$</field>
    <field name="win.eventdata.targetObject" type="pcre2">(?i)(System\\CurrentControlSet\\Control\\Terminal\sServer\\fDenyTSConnections|System\\CurrentControlSet\\Control\\Terminal\sServer\\WinStations\\RDP-Tcp\\PortNumber)</field>
    <description>RDP registry configuration modification (Defense Evasion/Lateral Movement)</description>
    <mitre>
      <id>T1112</id>
    </mitre>
    <group>PR.AA-05</group>
  </rule>

  <!-- Rule 100034: WMI Event Subscription Created -->
  <rule id="100034" level="12">
    <decoded_as>json</decoded_as>
    <field name="win.system.providerName">Microsoft-Windows-Sysmon</field>
    <field name="win.system.eventID">^19$|^20$|^21$</field>
    <field name="win.eventdata.filterName" type="pcre2">(?i)(__EventFilter|__EventConsumer|__FilterToConsumerBinding)</field>
    <description>WMI event filter/consumer activity (Persistence)</description>
    <mitre>
      <id>T1546.003</id>
    </mitre>
    <group>DE.AE-02</group>
  </rule>

</group>
