# rules/DE.CM/DE.CM-09_schtasks-creation.xml
<!--
  Rule: Scheduled task creation detected via Sysmon (Persistence)
  NIST CSF: DE.AE-02
  MITRE ATT&CK: T1053.005
  Severity: High
  Author: Subhash Bharadwaj
  Last Updated: 2026-05-23
-->

<group name="nist_csf, de.cm-09, mitre_attack,">

  <rule id="100029" level="8">
      <decoded_as>json</decoded_as>
      <field name="win.system.providerName">Microsoft-Windows-Sysmon</field>
      <field name="win.system.eventID">^1$</field>
      <field name="win.eventdata.commandLine" type="pcre2">(?i)schtasks\.exe.*\/create</field>
      <description>Scheduled task creation detected via Sysmon (Persistence)</description>
      <mitre>
        <id>T1053.005</id>
      </mitre>
      <group>DE.AE-02</group>
    </rule>

</group>
